# scripts/com.funtastic.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.funtastic.worker</string>

    <key>WorkingDirectory</key>
    <string>/Users/ian/Desktop/funtastic-saas</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>--env-file=.env.local</string>
        <string>--import</string>
        <string>tsx</string>
        <string>src/worker.ts</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ian/Desktop/funtastic-saas/logs/worker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ian/Desktop/funtastic-saas/logs/worker-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
